# target/linux/ipq60xx/base-files/etc/hotplug.d/ieee80211/10-fix-mac-address
[ "$ACTION" == "add" ] || exit 0

PHY_NUM=${DEVPATH##*/phy}

[ -n $PHY_NUM ] || exit 0

. /lib/functions.sh
. /lib/functions/system.sh

case "$(board_name)" in
esac
